# docker/frankenphp/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/www/storage/logs/supervisord_frankenphp.log
pidfile=/var/run/supervisord_frankenphp.pid

[program:octane-frankenphp]
command=php artisan octane:start --server=frankenphp --host=0.0.0.0 --port=80 --admin-port=2019
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:queue-worker]
command=php artisan queue:work --tries=3 --timeout=90
autostart=true
autorestart=true
numprocs=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:schedule-worker]
command=php artisan schedule:work
autostart=true
autorestart=true
numprocs=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
